# etc/cron.weekly/geoip
#!/bin/sh

test -x /usr/sbin/geoip-database-contrib_update || exit 0
/usr/sbin/geoip-database-contrib_update