# etc/cron.weekly/mysqloptimize
#!/bin/sh

test -x /usr/bin/mysqloptimize || exit 0

/usr/bin/mysqloptimize -A

test -x /usr/bin/mysqlrepair || exit 0

/usr/bin/mysqlrepair -A